# recipes-bsp/optee/files/thsm.service
[Unit]
Description=TEE-based HSM Service
After=optee.service

[Service]
ExecStart=/usr/bin/thsmDaemon

[Install]
WantedBy=multi-user.target
